# nms.conf
description "4ipnet nms manager demon"

start on filesystem
script
    echo 2097152 > /proc/sys/fs/nr_open
    ulimit -n 1048576
    /usr/bin/nms &
end script
