# sflow.conf
description "4ipnet sflow read demon"

start on filesystem
script
    echo 2097152 > /proc/sys/fs/nr_open
    ulimit -n 1048576
    /usr/bin/sflowtool &
end script
